# playbooks/data/teak-services.slice
[Unit]
Description=Slice for Teak support services

[Slice]
MemoryMax=96M
