# playbooks/data/teak-services.slice
[Unit]
Description=Slice for Teak support services

[Slice]
CPUQuota=20%
MemoryMax=64M
