# playbooks/data/update_logrotate.timer
[Timer]
OnCalendar=
OnCalendar=hourly
AccuracySec=1m
